# deploy/cybersentinels.service
# Systemd service for Cyber Sentinels (Daphne - HTTP + WebSockets)
# Copy to: sudo cp deploy/cybersentinels.service /etc/systemd/system/
# Then: Replace YOUR_USER and /home/YOUR_USER with your actual user and project path
#       sudo systemctl daemon-reload && sudo systemctl enable cybersentinels && sudo systemctl start cybersentinels

[Unit]
Description=Cyber Sentinels Website (Daphne ASGI)
After=network.target redis-server.service

[Service]
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/Cyber-Sentinels-Website
Environment="PATH=/home/YOUR_USER/Cyber-Sentinels-Website/venv/bin"
ExecStart=/home/YOUR_USER/Cyber-Sentinels-Website/venv/bin/daphne -b 127.0.0.1 -p 8000 Cybersentinels_website.asgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
